# systemd/scheduled-reboot.timer
[Unit]
Description=Weekly reboot of access control

[Timer]
OnCalendar=Thu 14:30:00

[Install]
WantedBy=timers.target
