# systemd/scheduled-reboot.service
[Unit]
Description=Daily Reoot of access control

[Service]
Type=simple
ExecStart=/usr/bin/systemctl --force reboot
